# dot_config/systemd/user/onedrive_mount.service
[Unit]
Description=Onedrive
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
ExecStart=/usr/bin/rclone mount onedrive: /home/shejialuo/Network/onedrive \
    --config=/home/shejialuo/.config/rclone/rclone.conf \
    --cache-dir=/home/shejialuo/.cache/rclone/vfs \
    --vfs-cache-mode=full \
    --vfs-read-chunk-streams=8 \
    --vfs-read-chunk-size=16M
ExecStop=/usr/bin/fusermount -u /home/shejialuo/Network/onedrive
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
